# deploy/systemd/minsy-flower.service
[Unit]
Description=Minsy Celery Flower Monitor
After=network.target

[Service]
Type=simple
User=yifan
WorkingDirectory=/home/yifan/minsy_mvp_backend_v2
EnvironmentFile=/home/yifan/minsy_mvp_backend_v2/.env
ExecStart=/usr/bin/env bash -lc 'cd /home/yifan/minsy_mvp_backend_v2 && exec uv run celery -A src.workers.celery_app.celery_app flower --address="${FLOWER_HOST:-127.0.0.1}" --port="${FLOWER_PORT:-5555}" --basic_auth="${FLOWER_USER}:${FLOWER_PASSWORD}"'
Restart=always
RestartSec=5
TimeoutStartSec=30

[Install]
WantedBy=multi-user.target

